# local/.config/autostart/ksuperkey.desktop
[Desktop Entry]
Exec=ksuperkey -t 200
Name=ksuperkey
Type=Application
Version=1.0
